# ansible/roles/log_server/files/remote-log-rotate.service
[Unit]
Description=Compress/cleanup aggregated remote logs in /var/log/remote

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/remote-log-rotate.sh
